# projects/system-health/vm-deploy/systemd/ynai5-gemini.service
[Unit]
Description=YNAI5 Gemini Worker
After=ynai5-drive.service network-online.target
Requires=ynai5-drive.service

[Service]
Type=simple
User=shema
WorkingDirectory=/ynai5_runtime
EnvironmentFile=/ynai5_runtime/.env
ExecStart=/usr/bin/python3 /ynai5_runtime/scripts/gemini_worker.py
Restart=always
RestartSec=10
StandardOutput=append:/ynai5_runtime/logs/gemini_worker.log
StandardError=append:/ynai5_runtime/logs/gemini_worker.log

[Install]
WantedBy=multi-user.target
